# deploy/lokaApps.service
[Unit]
Description=loka apps (Node API backend)
After=network.target

[Service]
Type=simple
User=mithun
WorkingDirectory=/home/mithun/loka.place/lokaApps/api
EnvironmentFile=/home/mithun/loka.place/lokaApps/api/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
